# sd/monarch_compress.timer
[Timer]
OnBootSec=0min
OnUnitActiveSec=1min
#OnCalendar=*:0/1

[Install]
WantedBy=timers.target
