# sd/monarch_arecord.service
[Service]
ExecStart=/data/p/monarch/bin/monarch arecord
